# backup.timer
[Unit]
Description=Run backup every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target